# resources/etc/services.d/varnish/run
#!/usr/bin/sh
/usr/sbin/varnishd -T localhost:6082 -a :80 -s default=malloc,256M -f /etc/varnish/default.vcl
echo "varnishd starting"